# services/com.clipshare.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clipshare.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/clip_client</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLIPSHARE_TOKEN</key>
        <string>your_token_here</string>
        <key>CLIPSHARE_POLL_INTERVAL</key>
        <string>2</string>
        <key>CLIPSHARE_VERBOSE</key>
        <string>false</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/clipshare.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/clipshare.error</string>

    <key>WorkingDirectory</key>
    <string>/tmp</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
